# transfer_funds.cron
* * * * * /home/ubuntu/creator-classes/creatorenv/bin/python3 /home/ubuntu/creator-classes/manage.py weekly_payouts >> /home/ubuntu/creator-classes/cron.log 2>&1
* * * * * /home/ubuntu/creator-classes/creatorenv/bin/python3 /home/ubuntu/creator-classes/manage.py session_notification >> /home/ubuntu/creator-classes/cron.log 2>&1
* * * * * /home/ubuntu/creator-classes/creatorenv/bin/python3 /home/ubuntu/creator-classes/manage.py stream_notification >> /home/ubuntu/creator-classes/cron.log 2>&1
